# packaging/deb/debian/prerm
#!/bin/sh -e

systemctl disable --now kat-rocky.socket
systemctl disable --now kat-rocky.service
