# packaging/deb/debian/prerm
#!/bin/sh -e

systemctl disable --now kat-mula.service
