# packaging/deb/debian/prerm
#!/bin/sh -e

systemctl disable --now kat-octopoes.socket
systemctl disable --now kat-octopoes.service
