# rocky/debian/prerm
#!/bin/sh -e

if [ "$1" != "upgrade" ]; then
    systemctl disable --now kat-rocky.service
fi

#DEBHELPER#
