# octopoes/debian/prerm
#!/bin/sh -e

if [ "$1" != "upgrade" ]; then
    systemctl disable --now kat-octopoes.service
fi

#DEBHELPER#
